# daily-photo-button.service
[Unit]
Description=Daily Photo Button Monitor
After=network.target
Wants=network.target

[Service]
User=pi
Group=pi
WorkingDirectory=/home/pi/daily-photo-client
ExecStart=/home/pi/daily-photo-client/run_button_monitor.sh
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/home/pi/daily-photo-client
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
